# conf/crontab
# Run `proxyctl cron` every minute
* * * * * /usr/local/bin/proxyctl cron
